# launchd/com.brickblade.refresh-prices.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Nightly price refresh for owned sets. Runs at 04:00 local time.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brickblade.refresh-prices</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>--project</string>
    <string>PROJECT_DIR</string>
    <string>brickblade</string>
    <string>refresh-prices</string>
  </array>

  <key>WorkingDirectory</key>
  <string>PROJECT_DIR</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>PROJECT_DIR/var/logs/refresh-prices.stdout</string>

  <key>StandardErrorPath</key>
  <string>PROJECT_DIR/var/logs/refresh-prices.stderr</string>
</dict>
</plist>
